# deploy/vigo.service
# ==============================================
# Vigo 应用 Systemd 服务配置文件
# 适用于 CentOS 7.x 及以上版本
# ==============================================
# 使用方法:
# 1. 将此文件复制到 /etc/systemd/system/vigo.service
# 2. 修改 WorkingDirectory 和 ExecStart 路径为实际部署路径
# 3. 执行: systemctl daemon-reload
# 4. 启动: systemctl start vigo
# 5. 开机自启: systemctl enable vigo
# ==============================================

[Unit]
Description=Vigo Web Application
Documentation=https://github.com/your-org/vigo
After=network.target mysql.service redis.service

[Service]
# 应用运行用户（建议使用非 root 用户）
User=www-data
Group=www-data

# 应用部署目录
WorkingDirectory=/opt/vigo

# 启动命令
ExecStart=/opt/vigo/vigo-linux-amd64

# 重启策略
Restart=always
RestartSec=5

# 环境变量（可选）
# Environment="DB_HOST=127.0.0.1"
# Environment="DB_PORT=3306"
# Environment="REDIS_HOST=127.0.0.1"

# 文件描述符限制（高并发优化）
LimitNOFILE=65536

# 标准输出和错误输出
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
